# bootc-upgrade-silent.service
[Unit]
Description=Apply bootc updates silently (no reboot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bootc upgrade